# meta-dpc/recipes-dpc/dpc-systemd-demo/dpc-systemd-demo/services/stress-cpu-starvation.service
[Unit]
Description=Starvate CPU

[Service]
ExecStart=/usr/bin/stress -c 1
CPUQuota=20%

[Install]
WantedBy=multi-user.target
